# etc/s6-overlay/s6-rc.d/clean-sed/up
#!/command/execlineb -P
foreground { echo "Calling /etc/rc-user/clean-sed..." }
/bin/bash /etc/rc-user/clean-sed
