# etc/s6-overlay/s6-rc.d/dot-upstream/up
#!/command/execlineb -P
foreground { echo "Calling /etc/rc-user/dot-upstream..." }
/bin/bash /etc/rc-user/dot-upstream